# files/etc/systemd/system/postgresql-delayed.timer
[Unit]
Description=Delayed PostgreSQL Service Timer

[Timer]
OnBootSec=60
Unit=postgresql-delayed.service

[Install]
WantedBy=timers.target
